# systemd/dot-config/systemd/user/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=30m
Unit=mbsync.service

[Install]
WantedBy=timers.target
